# wallpaper-shuffle.service
[Unit]
Description=Shuffles the current wallpaper

[Service]
Type=oneshot
ExecStart=%h/.local/bin/wallpaper-shuffle
